# AUR/rgb-cooling-hat.install
post_install() {
    systemctl enable --now i2c-fan
}

pre_remove() {
    systemctl disable --now i2c-fan
}
