# AUR/rgb-cooling-hat.install
post_install() {
    systemctl enable --now rgb-cooling-hat
}

pre_remove() {
    systemctl disable --now rgb-cooling-hat
}
